# contrib/packaging/freebsd/multimedia/mympd/pkg-plist
bin/mympd
bin/mympd-config
bin/mympd-script
etc/rc.d/mympd
share/man/man1/mympd.1.gz
share/man/man1/mympd-config.1.gz
share/man/man1/mympd-script.1.gz
