# contrib/packaging/freebsd/multimedia/mympd/pkg-plist
bin/mympd
etc/rc.d/mympd
share/man/man1/mympd.1.gz
