# contrib/initscripts/mympd.openrc.in
#!/sbin/openrc-run
#
# SPDX-License-Identifier: GPL-3.0-or-later
# myMPD (c) 2018-2022 Juergen Mang <mail@jcgames.de>
# https://github.com/jcorporation/mympd

name=myMPD
command="@CMAKE_INSTALL_FULL_BINDIR@/mympd"
command_args="-s"
command_background="yes"

start_stop_daemon_args=""
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
