# debian/postrm
if [ -f /usr/lib/systemd/system/mympd.service ]
then
  if `systemctl is-active --quiet mympd`
  then
    echo "stopping mympd.service" && systemctl stop mympd 
  fi
  echo "disabling mympd.service" && systemctl disable mympd
  rm -v -f /usr/lib/systemd/system/mympd.service
  systemctl daemon-reload
fi

rm -v -f /usr/share/mympd/htdocs/pics
rm -v -f /usr/share/mympd/htdocs/library
